# files/openchat.conf
[program:openchat]
command = /webapps/openchat_env/bin/gunicorn_start
user = openchat
stdout_logfile = /webapps/openchat_env/logs/gunicorn_supervisor.log
redirect_stderr = true
